# mpvpaper/mpvpaper-back.service

[Unit]
Description=mpvpaper for the back monitor
After=graphical-session.target wayland-session.target
Wants=wayland-session.target

[Service]
Type=simple
ExecStart=%h/.config/hypr/mpvpaper/mpvpaper-back.sh
Restart=on-failure
RestartSec=5
StandardOutput=null
StandardError=journal
PassEnvironment=WAYLAND_DISPLAY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS HYPRLAND_INSTANCE_SIGNATURE

[Install]
WantedBy=default.target

